# kubernetes/apps/network/omada-cert-sync/brain/omada-cert-sync.service
[Unit]
Description=Sync cluster LE wildcard cert into Omada controller keystore
Documentation=https://github.com/rwlove/home-ops/tree/main/kubernetes/apps/network/omada-cert-sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/omada-cert-sync.sh
# Hardening: the script needs root (keystore + tpeap), but lock down the rest.
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
