# kubernetes/apps/network/omada-cert-sync/brain/omada-cert-sync.timer
[Unit]
Description=Daily check for renewed cluster wildcard cert (Omada portal)

[Timer]
# Daily at 03:00 local — lands in the routine maintenance window so the
# rare renewal-day Omada restart (~quarterly) is in-window. Every other
# tick is a no-op fingerprint compare.
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
